# Scripts/Task_90Days.xml
<?xml version="1.0" encoding="UTF-16"?>

<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">

  <RegistrationInfo>

    <Date>2026-06-10T11:01:06.8869385</Date>

    <Author>DESKTOP-NT9JMGD\koriat</Author>

    <Description>SmartScale Worker — סנכרון 90 יום (פעם ביום בשעה 18:00)</Description>

    <URI>\SmartScaleWorkerTask_90Days</URI>

  </RegistrationInfo>

  <Principals>

    <Principal id="Author">

      <UserId>S-1-5-21-3543236145-1136098358-3405112608-1001</UserId>

      <LogonType>Password</LogonType>

      <RunLevel>HighestAvailable</RunLevel>

    </Principal>

  </Principals>

  <Settings>

    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>

    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>

    <ExecutionTimeLimit>PT15M</ExecutionTimeLimit>

    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>

    <StartWhenAvailable>true</StartWhenAvailable>

    <IdleSettings>

      <StopOnIdleEnd>true</StopOnIdleEnd>

      <RestartOnIdle>false</RestartOnIdle>

    </IdleSettings>

  </Settings>

  <Triggers>

    <CalendarTrigger>

      <StartBoundary>2026-06-10T18:00:00</StartBoundary>

      <ScheduleByDay>

        <DaysInterval>1</DaysInterval>

      </ScheduleByDay>

    </CalendarTrigger>

  </Triggers>

  <Actions Context="Author">

    <Exec>

      <Command>C:\inetpub\wwwroot\migrateDB\Worker\RunWorker_90Days.bat</Command>

      <WorkingDirectory>C:\inetpub\wwwroot\migrateDB\Worker</WorkingDirectory>

    </Exec>

  </Actions>

</Task>
